# etc/services.d/earnapp/run
#!/command/with-contenv bash

exec /usr/bin/earnapp_service
sleep 5
